# ntfs-automount/ntfs-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/bin/ntfs-automount %I
